# devops/rpmbuild/extrafiles/DEBIAN-ethermint/postinst
#!/bin/sh
# postinst script for basecoin
#

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    chown ethermint.ethermint /etc/ethermint
    chown ethermint.ethermint /etc/ethermint/genesis.json
    chown ethermint.ethermint /etc/ethermint/tendermint
    chown ethermint.ethermint /etc/ethermint/tendermint/config.toml
    chown ethermint.ethermint /etc/ethermint/keystore
    chown ethermint.ethermint /etc/ethermint/keystore/UTC--2016-10-21T22-30-03.071787745Z--7eff122b94897ea5b0e2a9abf47b86337fafebdc

    /usr/bin/ethermint --datadir /etc/ethermint init /etc/ethermint/genesis.json
    chown -R ethermint.ethermint /etc/ethermint/ethermint
    test ! -f /etc/ethermint/tendermint/priv_validator.json && tendermint gen_validator > /etc/ethermint/tendermint/priv_validator.json && chmod 0400 /etc/ethermint/tendermint/priv_validator.json && chown ethermint.ethermint /etc/ethermint/tendermint/priv_validator.json
    tendermint_pubkey="`tendermint show_validator --home /etc/ethermint/tendermint --log_level error`"
    chown ethermint.ethermint /etc/ethermint/tendermint/data
    test ! -f /etc/ethermint/tendermint/genesis.json && cat << EOF > /etc/ethermint/tendermint/genesis.json
{
  "genesis_time": "2017-06-10T03:37:03Z",
  "chain_id": "my_chain_id",
  "validators":
  [
    {
      "pub_key": $tendermint_pubkey,
      "amount":10,
      "name":"my_testchain_node"
    }
  ],
  "app_hash": "",
  "app_options": {}
}
EOF
    chown ethermint.ethermint /etc/ethermint/tendermint/genesis.json
    systemctl daemon-reload
    systemctl enable ethermint
    if [ -d /etc/ethermint/tendermint/data ]; then
      service ethermint start
    fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
