# tools/build/gaia/etc/systemd/system/gaia.service
[Unit]
Description=Gaia
Requires=network-online.target
After=network-online.target

[Service]
Restart=on-failure
User=gaia
Group=gaia
PermissionsStartOnly=true
ExecStart=/usr/bin/gaia node start --home=/etc/gaia
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

