# build/gaia/etc/systemd/system/gaia-server.service
[Unit]
Description=Gaia server
Requires=network-online.target
BindTo=gaia.service
PartOf=gaia.service
Before=gaia.service
After=network-online.target
PropagatesReloadTo=gaia.service
ReloadPropagatedFrom=gaia.service

[Service]
Environment="TMHOME=/etc/gaia/tendermint"
Restart=on-failure
User=gaia
Group=gaia
PermissionsStartOnly=true
ExecStart=/usr/bin/tendermint node
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
Also=gaia.service

